# odc_gee/opt/config/systemd/user/update_products.timer
[Unit]
Description=Update products daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
